# BlogBuild.service
[Unit]
Description=Build Blog.

[Service]
User=User
Type=oneshot
ExecStart=path-to-script

[Install]
WantedBy=multi-user.target

